# configs/fireadmin_auth.conf
[program:fireadmin_auth]
command=python3.6 main.py
directory = /root/fireadmin/auth
autostart = true
stderr_logfile = /var/log/fireadmin_auth.log
stdout_logfile = /var/log/fireadmin_auth.log