# configs/fireadmin.conf
[program:fireadmin]
command=python3.6 main.py
directory = /root/fireadmin/app
autostart = true
stderr_logfile = /var/log/fireadmin.log
stdout_logfile = /var/log/fireadmin.log